# src/batterylevel.service
[Unit]
Description=batterylevel

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/batterylevel -m
